# service/vmimage/overrides/etc/periodic/daily/update
#!/usr/bin/env sh
apk update
apk upgrade
